# supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/data/supervisord.pid

[program:runner_manager]
command=python -m app.runner_manager
directory=/app
autostart=true
autorestart=true
stdout_logfile=/data/logs/manager.log
stderr_logfile=/data/logs/manager.log

[program:web]
command=gunicorn -w 1 -b 0.0.0.0:8080 --access-logfile - --error-logfile - app.web:app
directory=/app
autostart=true
autorestart=true
stdout_logfile=/data/logs/web.log
stderr_logfile=/data/logs/web.log
